# install/systemd/radiopi.service
[Unit]
Description=Radiopi

[Service]
Type=simple
ExecStart=/usr/bin/radiopi

[Install]
WantedBy=multi-user.target
